# ansible/roles/transcoder/templates/systemd-units/transcode.target.j2
#MANAGED BY ANSIBLE

[Unit]
Description = All Transcoders of Stream {{ item }}

[Install]
WantedBy = multi-user.target
