# ansible/roles/transcoder/templates/systemd-units/transcode.service.j2
#MANAGED BY ANSIBLE

[Unit]
Description = Transcoder {{ item.1 }} Stream s{{ item.0 }}
Requires = transcode_s{{ item.0 }}.target
After = transcode_s{{ item.0 }}.target

[Service]
Type = simple
ExecStart = /opt/transcoder/scripts/s{{ item.0 }}_{{ item.1 }}.sh
Restart = always
RestartSec = 1s
StartLimitInterval = 0

[Install]
WantedBy = transcode_s{{ item.0 }}.target
