# ansible/roles/transcoder/templates/systemd-units/transcode.service.j2
#MANAGED BY ANSIBLE

[Unit]
Description = Transcode {{ item }} Stream %I
Requires = transcode@%i.target
After = transcode@%i.target

[Service]
Type = simple
EnvironmentFile = /opt/transcoder/config/%i
ExecStart = /opt/transcoder/scripts/{{ item }}.sh
Restart = always
RestartSec = 5s
StartLimitInterval = 0
