# ansible/roles/transcoder/templates/systemd-units/transcode.target.j2
#MANAGED BY ANSIBLE

[Unit]
Description = All Transcoders of Stream %I
Wants = transcode@%i.service

[Install]
WantedBy = multi-user.target
