# ansible/roles/transcoder/templates/systemd-units/update_transcoding.service.j2
#MANAGED BY ANSIBLE

[Unit]
Description = Update transcoding tasks

[Service]
Type = simple
ExecStart = /opt/transcoder/scripts/update_transcoding.py -u {{ stream_api.user }} -p {{ stream_api.password }} {{ stream_api.url }} {{ inventory_hostname }}
StartLimitInterval = 0
