# ansible/roles/audio-pi/files/audio-ingest-failure-notice@.service
[Unit]
Description = audio-ingest-failure-notice for %i

[Service]
Type = oneshot
ExecStart = /opt/mqtt/alert-failed.sh

