# ansible/roles/relay/templates/systemd-units/fanout.target.j2
#MANAGED BY ANSIBLE

[Unit]
Description = All Fanouts of Stream s{{ item }}

[Install]
WantedBy = multi-user.target
