# ansible/roles/relay/templates/systemd-units/fanout_dash.service.j2
#MANAGED BY ANSIBLE

[Unit]
Description = Fanout Dash Stream %I
Requires = fanout@%i.target
After = fanout@%i.target

[Service]
Type = simple
ExecStart = /opt/fanout/scripts/dash.py \
	--stream=%i \
	--pull_endpoint={{ relay_fanout.pull_endpoint }} \
	--dash_write_path={{ relay_fanout.dash_write_path }}

Restart = always
RestartSec = 30s
StartLimitInterval = 0

[Install]
WantedBy = fanout@%i.target
