# ansible/roles/relay/templates/systemd-units/stream_update.timer.j2
[Unit]
Description=Update stream info every minute

[Timer]
OnCalendar=*-*-* *:*:00/30
AccuracySec=10s

[Install]
WantedBy=timers.target
