# ansible/roles/relay/templates/systemd-units/update_fanout.timer.j2
[Unit]
Description=Update fanout every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
