# bundlewrap/bundles/zfs/files/zfs-import-scan-override.service
[Unit]
ConditionPathExists=

[Service]
ExecStart=
ExecStart=/usr/sbin/zpool import -aN -o cachefile=none
