# ansible/roles/letsencrypt/files/dehydrated/dehydrated.service
[Unit]
Description=Run dehydrated to renew Let's Encrypt certificates

[Service]
Type=oneshot
ExecStart=/opt/dehydrated/dehydrated --cron
