# ansible/roles/mixer-obs/templates/systemd-units/x11vnc-mouseable.service.j2
[Unit]
Description=x11vnc

[Service]
ExecStart=/usr/bin/x11vnc -display :0 -alwaysshared -dontdisconnect -shared -forever -listen 127.0.0.1 -rfbport 5901 -no6 -noipv6 -listen6 ""
User=voc
Group=voc

[Install]
WantedBy=multi-user.target
