# ansible/roles/loadbalancer/templates/mqttfeedback-restart.path.j2
[Unit]
Description=Mqttfeedback File Watch

[Path]
PathChanged=/srv/nginx/streaming-website/lib/feedback/feedback.py

[Install]
WantedBy=multi-user.target