# ansible/roles/loadbalancer/templates/schedule-download.service.j2
[Unit]
Description=Download schedules for streaming-website

[Service]
Type=oneshot
ExecStart=sudo -udownloader /usr/bin/env php /srv/nginx/streaming-website/index.php download