# ansible/roles/loadbalancer/templates/schedule-download.timer.j2
[Unit]
Description=Update schedules every 15min

[Timer]
OnBootSec=1min
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
