# ansible/roles/relay/files/stats/stats.service
[Unit]
Description=Player stats server

[Service]
ExecStart=/usr/local/bin/voc-telemetry
Restart=on-failure

[Install]
WantedBy=multi-user.target
