# bundlewrap/bundles/crs-worker/files/crs-worker.target
[Unit]
Description=CRS worker scripts
Requires=network.target
