# ansible/roles/relay/templates/register_relay/register_relay.service.j2
#MANAGED BY ANSIBLE

[Unit]
Description = Register Relay service

[Service]
ExecStart = /usr/local/sbin/register_relay
Type = oneshot
RemainAfterExit = yes

[Install]
WantedBy = multi-user.target
